# config/uniopen.conf
[program:uniopen-backend]
command = gunicorn manage:app -b localhost:5000
directory = /home/uniopen/uniopen
user = uniopen