# config/unipd.conf
[program:uniopen-backend]
command = gunicorn app:app -b localhost:8000
directory = /home/uniopen/bot_app
user = uniopen