# config/unipd.conf
[program:unipd]
command = gunicorn app:app -b localhost:8000
directory = /home/mikexine/bot_app
user = mikexine